# wsl-secret-service.service
[Unit]
Description=WSL2 Secret Service (Windows Credential Manager backend)
Documentation=https://specifications.freedesktop.org/secret-service/0.2/
# Start after the user's D-Bus session bus is available.
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart=%h/.local/bin/wsl-secret-service
Restart=on-failure
RestartSec=3

# WSL2 sets DBUS_SESSION_BUS_ADDRESS via /run/user/<UID>/bus when
# systemd user instance is active.  Override here if needed.
# Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=default.target
